# deploy/voltra.service
# systemd unit: start the Voltra stack on boot (Linux VPS).
# Install:
#   sudo cp deploy/voltra.service /etc/systemd/system/voltra.service
#   sudo systemctl daemon-reload && sudo systemctl enable --now voltra
# Assumes the repo lives at /opt/voltra (edit WorkingDirectory if not).

[Unit]
Description=Voltra trading stack (Docker Compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/voltra
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.prod.yml down
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
